# yhafri/rabbitmq/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/rabbitmqd
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/ant/
drwxr-xr-x	root/root	usr/lib/ant/lib/
-rw-r--r--	root/root	usr/lib/ant/lib/rabbitmq-client-tests.jar
-rw-r--r--	root/root	usr/lib/ant/lib/rabbitmq-client.jar
drwxr-xr-x	root/root	usr/lib/erlang/
drwxr-xr-x	root/root	usr/lib/erlang/lib/
drwxr-xr-x	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/
drwxr-xr-x	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit.app
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit.rel
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_access_control.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_amqqueue.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_amqqueue_process.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_binary_generator.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_binary_parser.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_channel.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_channel0.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_control.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_error_logger.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_exchange.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_framing.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_gc_persist.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_gensym.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_http.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_http_channel.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_log.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_management.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_misc.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_persister.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_reader.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_realm.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_tests.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_ticket.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_tracer.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_transaction.beam
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/ebin/rabbit_writer.beam
drwxr-xr-x	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/include/
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/include/rabbit.hrl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/include/rabbit_framing.hrl
drwxr-xr-x	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_access_control.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_amqqueue.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_amqqueue_process.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_binary_generator.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_binary_parser.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_channel.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_channel0.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_control.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_error_logger.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_exchange.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_framing.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_gc_persist.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_gensym.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_http.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_http_channel.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_log.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_management.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_misc.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_persister.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_reader.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_realm.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_tests.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_ticket.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_tracer.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_transaction.erl
-rw-r--r--	root/root	usr/lib/erlang/lib/rabbitmq-1.0.0-alpha/src/rabbit_writer.erl
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/rabbitmq-server
-rw-r--r--	root/root	usr/sbin/rabbitmq-server.bat
-rw-r--r--	root/root	usr/sbin/rabbitmqctl
-rw-r--r--	root/root	usr/sbin/rabbitmqctl.bat
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/rabbitmq/
drwxr-xr-x	root/root	var/lib/rabbitmq/mnesia/
drwxr-xr-x	root/root	var/log/
drwxr-xr-x	root/root	var/log/rabbitmq/
