# yhafri/heritrix/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/heritrix
drwxrwxrwt	root/root	tmp/
drwxr-xr-x	root/root	tmp/heritrix/
drwxr-xr-x	root/root	tmp/heritrix/jobs/
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/heritrix-1.12.1/
drwxr-xr-x	root/root	usr/lib/heritrix-1.12.1/bin/
-rwxr-xr-x	root/root	usr/lib/heritrix-1.12.1/bin/arcreader
-rwxr-xr-x	root/root	usr/lib/heritrix-1.12.1/bin/arcreader.cmd
-rwxr-xr-x	root/root	usr/lib/heritrix-1.12.1/bin/cmdline-jmxclient-0.10.5.jar
-rwxr-xr-x	root/root	usr/lib/heritrix-1.12.1/bin/extractor
-rwxr-xr-x	root/root	usr/lib/heritrix-1.12.1/bin/extractor.cmd
-rwxr-xr-x	root/root	usr/lib/heritrix-1.12.1/bin/foreground_heritrix
-rwxr-xr-x	root/root	usr/lib/heritrix-1.12.1/bin/foreground_heritrix.cmd
-rwxr-xr-x	root/root	usr/lib/heritrix-1.12.1/bin/heritrix
-rwxr-xr-x	root/root	usr/lib/heritrix-1.12.1/bin/heritrix.cmd
-rwxr-xr-x	root/root	usr/lib/heritrix-1.12.1/bin/hoppath.pl
-rwxr-xr-x	root/root	usr/lib/heritrix-1.12.1/bin/htmlextractor
-rwxr-xr-x	root/root	usr/lib/heritrix-1.12.1/bin/htmlextractor.cmd
-rwxr-xr-x	root/root	usr/lib/heritrix-1.12.1/bin/make_reports.pl
drwxr-xr-x	root/root	usr/lib/heritrix-1.12.1/conf/
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/conf/heritrix.cacerts
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/conf/heritrix.properties
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/conf/jmxremote.password.template
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/heritrix-1.12.1.jar
drwxr-xr-x	root/root	usr/lib/heritrix-1.12.1/lib/
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/ant-1.6.2.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/bsh-2.0b4.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/commons-cli-1.0.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/commons-codec-1.3.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/commons-collections-3.1.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/commons-httpclient-3.0.1.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/commons-lang-2.3.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/commons-logging-1.0.4.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/commons-net-1.4.1.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/commons-pool-1.3.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/dnsjava-2.0.3.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/fastutil-5.0.3-heritrix-subset-1.0.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/itext-1.2.0.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/jasper-compiler-tomcat-4.1.30.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/jasper-runtime-tomcat-4.1.30.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/javaswf-CVS-SNAPSHOT-1.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/je-3.2.23.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/jericho-html-2.3.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/jets3t-0.5.0.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/jetty-4.2.23.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/junit-3.8.2.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/libidn-0.5.9.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/mg4j-1.0.1.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/poi-2.0-RC1-20031102.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/poi-scratchpad-2.0-RC1-20031102.jar
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/lib/servlet-tomcat-4.1.30.jar
drwxr-xr-x	root/root	usr/lib/heritrix-1.12.1/webapps/
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/webapps/admin.war
-rw-r--r--	root/root	usr/lib/heritrix-1.12.1/webapps/selftest.war
