# ports/fail2ban/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/fail2ban/
drwxr-xr-x	root/root	etc/fail2ban/action.d/
-rw-rw-r--	root/root	etc/fail2ban/action.d/complain.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/dshield.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/dummy.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/hostsdeny.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/ipfilter.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/ipfw.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/iptables-allports.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/iptables-multiport-log.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/iptables-multiport.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/iptables-new.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/iptables.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/mail-buffered.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/mail-whois-lines.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/mail-whois.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/mail.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/mynetwatchman.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/sendmail-buffered.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/sendmail-whois-lines.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/sendmail-whois.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/sendmail.conf
-rw-rw-r--	root/root	etc/fail2ban/action.d/shorewall.conf
-rw-rw-r--	root/root	etc/fail2ban/fail2ban.conf
drwxr-xr-x	root/root	etc/fail2ban/filter.d/
-rw-rw-r--	root/root	etc/fail2ban/filter.d/apache-auth.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/apache-badbots.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/apache-nohome.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/apache-noscript.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/apache-overflows.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/common.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/courierlogin.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/couriersmtp.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/cyrus-imap.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/dovecot.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/dropbear.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/exim.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/gssftpd.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/lighttpd-fastcgi.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/named-refused.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/pam-generic.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/php-url-fopen.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/postfix.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/proftpd.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/pure-ftpd.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/qmail.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/sasl.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/sieve.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/sshd-ddos.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/sshd.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/vsftpd.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/webmin-auth.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/wuftpd.conf
-rw-rw-r--	root/root	etc/fail2ban/filter.d/xinetd-fail.conf
-rw-rw-r--	root/root	etc/fail2ban/jail.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-x---	root/root	etc/rc.d/fail2ban
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/fail2ban-client
-rwxr-xr-x	root/root	usr/bin/fail2ban-regex
-rwxr-xr-x	root/root	usr/bin/fail2ban-server
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/fail2ban/
drwxr-xr-x	root/root	usr/share/fail2ban/client/
-rw-r--r--	root/root	usr/share/fail2ban/client/__init__.py
-rw-r--r--	root/root	usr/share/fail2ban/client/__init__.pyc
-rw-r--r--	root/root	usr/share/fail2ban/client/actionreader.py
-rw-r--r--	root/root	usr/share/fail2ban/client/actionreader.pyc
-rw-r--r--	root/root	usr/share/fail2ban/client/beautifier.py
-rw-r--r--	root/root	usr/share/fail2ban/client/beautifier.pyc
-rw-r--r--	root/root	usr/share/fail2ban/client/configparserinc.py
-rw-r--r--	root/root	usr/share/fail2ban/client/configparserinc.pyc
-rw-r--r--	root/root	usr/share/fail2ban/client/configreader.py
-rw-r--r--	root/root	usr/share/fail2ban/client/configreader.pyc
-rw-r--r--	root/root	usr/share/fail2ban/client/configurator.py
-rw-r--r--	root/root	usr/share/fail2ban/client/configurator.pyc
-rw-r--r--	root/root	usr/share/fail2ban/client/csocket.py
-rw-r--r--	root/root	usr/share/fail2ban/client/csocket.pyc
-rw-r--r--	root/root	usr/share/fail2ban/client/fail2banreader.py
-rw-r--r--	root/root	usr/share/fail2ban/client/fail2banreader.pyc
-rw-r--r--	root/root	usr/share/fail2ban/client/filterreader.py
-rw-r--r--	root/root	usr/share/fail2ban/client/filterreader.pyc
-rw-r--r--	root/root	usr/share/fail2ban/client/jailreader.py
-rw-r--r--	root/root	usr/share/fail2ban/client/jailreader.pyc
-rw-r--r--	root/root	usr/share/fail2ban/client/jailsreader.py
-rw-r--r--	root/root	usr/share/fail2ban/client/jailsreader.pyc
drwxr-xr-x	root/root	usr/share/fail2ban/common/
-rw-r--r--	root/root	usr/share/fail2ban/common/__init__.py
-rw-r--r--	root/root	usr/share/fail2ban/common/__init__.pyc
-rw-r--r--	root/root	usr/share/fail2ban/common/helpers.py
-rw-r--r--	root/root	usr/share/fail2ban/common/helpers.pyc
-rw-r--r--	root/root	usr/share/fail2ban/common/protocol.py
-rw-r--r--	root/root	usr/share/fail2ban/common/protocol.pyc
-rw-r--r--	root/root	usr/share/fail2ban/common/version.py
-rw-r--r--	root/root	usr/share/fail2ban/common/version.pyc
drwxr-xr-x	root/root	usr/share/fail2ban/server/
-rw-r--r--	root/root	usr/share/fail2ban/server/__init__.py
-rw-r--r--	root/root	usr/share/fail2ban/server/__init__.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/action.py
-rw-r--r--	root/root	usr/share/fail2ban/server/action.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/actions.py
-rw-r--r--	root/root	usr/share/fail2ban/server/actions.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/asyncserver.py
-rw-r--r--	root/root	usr/share/fail2ban/server/asyncserver.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/banmanager.py
-rw-r--r--	root/root	usr/share/fail2ban/server/banmanager.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/datedetector.py
-rw-r--r--	root/root	usr/share/fail2ban/server/datedetector.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/datetemplate.py
-rw-r--r--	root/root	usr/share/fail2ban/server/datetemplate.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/faildata.py
-rw-r--r--	root/root	usr/share/fail2ban/server/faildata.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/failmanager.py
-rw-r--r--	root/root	usr/share/fail2ban/server/failmanager.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/failregex.py
-rw-r--r--	root/root	usr/share/fail2ban/server/failregex.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/filter.py
-rw-r--r--	root/root	usr/share/fail2ban/server/filter.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/filtergamin.py
-rw-r--r--	root/root	usr/share/fail2ban/server/filtergamin.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/filterpoll.py
-rw-r--r--	root/root	usr/share/fail2ban/server/filterpoll.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/iso8601.py
-rw-r--r--	root/root	usr/share/fail2ban/server/iso8601.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/jail.py
-rw-r--r--	root/root	usr/share/fail2ban/server/jail.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/jails.py
-rw-r--r--	root/root	usr/share/fail2ban/server/jails.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/jailthread.py
-rw-r--r--	root/root	usr/share/fail2ban/server/jailthread.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/mytime.py
-rw-r--r--	root/root	usr/share/fail2ban/server/mytime.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/server.py
-rw-r--r--	root/root	usr/share/fail2ban/server/server.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/ticket.py
-rw-r--r--	root/root	usr/share/fail2ban/server/ticket.pyc
-rw-r--r--	root/root	usr/share/fail2ban/server/transmitter.py
-rw-r--r--	root/root	usr/share/fail2ban/server/transmitter.pyc
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/run/
drwxr-xr-x	root/root	var/run/fail2ban/
