# ports/tor/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-x---	root/root	etc/rc.d/tor
drwxr-xr-x	root/root	etc/tor/
-rw-r-----	root/tor	etc/tor/torrc
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-x---	root/tor	usr/bin/tor
-rwxr-x---	root/tor	usr/bin/tor-gencert
-rwxr-x---	root/tor	usr/bin/tor-print-ed-signing-cert
-rwxr-x---	root/tor	usr/bin/tor-resolve
-rwxr-x---	root/tor	usr/bin/torify
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/tor-gencert.1.gz
-rw-r--r--	root/root	usr/share/man/man1/tor-print-ed-signing-cert.1.gz
-rw-r--r--	root/root	usr/share/man/man1/tor-resolve.1.gz
-rw-r--r--	root/root	usr/share/man/man1/tor.1.gz
-rw-r--r--	root/root	usr/share/man/man1/torify.1.gz
drwxr-xr-x	root/root	usr/share/tor/
-rw-r--r--	root/root	usr/share/tor/geoip
-rw-r--r--	root/root	usr/share/tor/geoip6
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	tor/tor	var/lib/tor/
-rw-------	tor/tor	var/lib/tor/cached-certs (EMPTY)
-rw-------	tor/tor	var/lib/tor/cached-consensus (EMPTY)
-rw-------	tor/tor	var/lib/tor/cached-descriptors (EMPTY)
-rw-------	tor/tor	var/lib/tor/cached-descriptors.new (EMPTY)
-rw-------	tor/tor	var/lib/tor/cached-routers.new (EMPTY)
-rw-------	tor/tor	var/lib/tor/state (EMPTY)
drwxr-xr-x	root/root	var/log/
drwxr-xr-x	tor/tor	var/log/tor/
-rw-r-----	tor/tor	var/log/tor/notices.log (EMPTY)
