# attic/lirc/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/lircd
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/modules/
drwxr-xr-x	root/root	lib/modules/<kernel-version>/
drwxr-xr-x	root/root	lib/modules/<kernel-version>/misc/
-rw-r--r--	root/root	lib/modules/<kernel-version>/misc/lirc_dev.ko
-rw-r--r--	root/root	lib/modules/<kernel-version>/misc/lirc_serial.ko
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/ircat
-rwxr-xr-x	root/root	usr/bin/irexec
-rwxr-xr-x	root/root	usr/bin/irpty
-rwxr-xr-x	root/root	usr/bin/irrecord
-rwxr-xr-x	root/root	usr/bin/irsend
-rwxr-xr-x	root/root	usr/bin/irw
-rwxr-xr-x	root/root	usr/bin/irxevent
-rwxr-xr-x	root/root	usr/bin/lircrcd
-rwxr-xr-x	root/root	usr/bin/mode2
-rwxr-xr-x	root/root	usr/bin/xmode2
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/lirc/
-rw-r--r--	root/root	usr/include/lirc/lirc_client.h
drwxr-xr-x	root/root	usr/lib/
-rw-r--r--	root/root	usr/lib/liblirc_client.a
-rwxr-xr-x	root/root	usr/lib/liblirc_client.la
lrwxrwxrwx	root/root	usr/lib/liblirc_client.so -> liblirc_client.so.0.2.0
lrwxrwxrwx	root/root	usr/lib/liblirc_client.so.0 -> liblirc_client.so.0.2.0
-rwxr-xr-x	root/root	usr/lib/liblirc_client.so.0.2.0
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/ircat.1.gz
-rw-r--r--	root/root	usr/man/man1/irexec.1.gz
-rw-r--r--	root/root	usr/man/man1/irpty.1.gz
-rw-r--r--	root/root	usr/man/man1/irrecord.1.gz
-rw-r--r--	root/root	usr/man/man1/irsend.1.gz
-rw-r--r--	root/root	usr/man/man1/irw.1.gz
-rw-r--r--	root/root	usr/man/man1/irxevent.1.gz
-rw-r--r--	root/root	usr/man/man1/lircrcd.1.gz
-rw-r--r--	root/root	usr/man/man1/mode2.1.gz
-rw-r--r--	root/root	usr/man/man1/smode2.1.gz
-rw-r--r--	root/root	usr/man/man1/xmode2.1.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/lircd.8.gz
-rw-r--r--	root/root	usr/man/man8/lircmd.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/lircd
-rwxr-xr-x	root/root	usr/sbin/lircmd
