# ports/conntrack-tools/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/conntrackd/
-rw-r--r--	root/root	etc/conntrackd/conntrackd.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/conntrackd
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/conntrack.8.gz
-rw-r--r--	root/root	usr/man/man8/conntrackd.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/conntrack
-rwxr-xr-x	root/root	usr/sbin/conntrackd
