# yhafri/tux/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/logrotate.d/
-rw-r--r--	root/root	etc/logrotate.d/tux
drwxr-xr-x	root/root	etc/rc.d/
drwxr-xr-x	root/root	etc/rc.d/init.d/
-rwxr-xr-x	root/root	etc/rc.d/init.d/tux
drwxr-xr-x	root/root	etc/sysconfig/
-rw-r--r--	root/root	etc/sysconfig/tux
-rw-r--r--	root/root	etc/tux.mime.types
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/include/
-rw-r--r--	root/root	usr/include/tux.h
-rw-r--r--	root/root	usr/include/tuxmodule.h
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/tux
-rwxr-xr-x	root/root	usr/sbin/tux2w3c
-rwxr-xr-x	root/root	usr/sbin/tuxstat
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man2/
-rw-r--r--	root/root	usr/share/man/man2/tux.2.gz
drwxr-xr-x	root/root	usr/share/man/man8/
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/tux/
