# yhafri/beagle/.footprint
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/beagle-config
-rwxr-xr-x	root/root	usr/bin/beagle-exercise-file-system
-rwxr-xr-x	root/root	usr/bin/beagle-extract-content
-rwxr-xr-x	root/root	usr/bin/beagle-index-info
-rwxr-xr-x	root/root	usr/bin/beagle-index-url
-rwxr-xr-x	root/root	usr/bin/beagle-info
-rwxr-xr-x	root/root	usr/bin/beagle-ping
-rwxr-xr-x	root/root	usr/bin/beagle-query
-rwxr-xr-x	root/root	usr/bin/beagle-shutdown
-rwxr-xr-x	root/root	usr/bin/beagle-status
-rwxr-xr-x	root/root	usr/bin/beagled
drwxr-xr-x	root/root	usr/etc/
drwxr-xr-x	root/root	usr/etc/beagle/
-rw-r--r--	root/root	usr/etc/beagle/crawl-applications
-rw-r--r--	root/root	usr/etc/beagle/crawl-documentation
-rw-r--r--	root/root	usr/etc/beagle/crawl-windows
-rw-r--r--	root/root	usr/etc/beagle/external-filters.xml
drwxr-xr-x	root/root	usr/etc/cron.daily/
-rwxr-xr-x	root/root	usr/etc/cron.daily/beagle-crawl-system
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/libbeagle/
drwxr-xr-x	root/root	usr/include/libbeagle/beagle/
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-client.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-daemon-information-request.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-daemon-information-response.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-empty-response.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-error-response.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-finished-response.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-hit.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-hits-added-response.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-hits-subtracted-response.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-indexable.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-indexing-service-request.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-property.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-query-part-date.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-query-part-human.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-query-part-or.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-query-part-property.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-query-part-text.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-query-part.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-query.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-request.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-response.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-search-term-response.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-shutdown-request.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-snippet-request.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-snippet-response.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-timestamp.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle-util.h
-rw-r--r--	root/root	usr/include/libbeagle/beagle/beagle.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/beagle/
drwxr-xr-x	root/root	usr/lib/beagle/Backends/
-rw-r--r--	root/root	usr/lib/beagle/Beagle.dll
-rw-r--r--	root/root	usr/lib/beagle/Beagle.dll.mdb
-rw-r--r--	root/root	usr/lib/beagle/BeagleDaemon.exe
-rw-r--r--	root/root	usr/lib/beagle/BeagleDaemon.exe.mdb
-rw-r--r--	root/root	usr/lib/beagle/BeagleDaemonLib.dll
-rw-r--r--	root/root	usr/lib/beagle/BeagleDaemonLib.dll.mdb
-rw-r--r--	root/root	usr/lib/beagle/BeagleDaemonPlugins.dll
-rw-r--r--	root/root	usr/lib/beagle/BeagleDaemonPlugins.dll.mdb
-rw-r--r--	root/root	usr/lib/beagle/BuildIndex.exe
-rw-r--r--	root/root	usr/lib/beagle/BuildIndex.exe.mdb
-rw-r--r--	root/root	usr/lib/beagle/Config.exe
-rw-r--r--	root/root	usr/lib/beagle/Config.exe.mdb
-rw-r--r--	root/root	usr/lib/beagle/ExerciseFileSystem.exe
-rw-r--r--	root/root	usr/lib/beagle/ExerciseFileSystem.exe.mdb
-rw-r--r--	root/root	usr/lib/beagle/ExtractContent.exe
-rw-r--r--	root/root	usr/lib/beagle/ExtractContent.exe.mdb
drwxr-xr-x	root/root	usr/lib/beagle/Filters/
-rw-r--r--	root/root	usr/lib/beagle/Filters/Filters.dll
-rw-r--r--	root/root	usr/lib/beagle/Filters/Filters.dll.mdb
-rw-r--r--	root/root	usr/lib/beagle/IndexHelper.exe
-rw-r--r--	root/root	usr/lib/beagle/IndexHelper.exe.mdb
-rw-r--r--	root/root	usr/lib/beagle/IndexWebContent.exe
-rw-r--r--	root/root	usr/lib/beagle/IndexWebContent.exe.mdb
-rw-r--r--	root/root	usr/lib/beagle/Info.exe
-rw-r--r--	root/root	usr/lib/beagle/Info.exe.mdb
-rw-r--r--	root/root	usr/lib/beagle/ManageIndex.exe
-rw-r--r--	root/root	usr/lib/beagle/ManageIndex.exe.mdb
-rw-r--r--	root/root	usr/lib/beagle/Query.exe
-rw-r--r--	root/root	usr/lib/beagle/Query.exe.mdb
-rw-r--r--	root/root	usr/lib/beagle/Shutdown.exe
-rw-r--r--	root/root	usr/lib/beagle/Shutdown.exe.mdb
-rw-r--r--	root/root	usr/lib/beagle/Util.dll
-rw-r--r--	root/root	usr/lib/beagle/Util.dll.config
-rw-r--r--	root/root	usr/lib/beagle/Util.dll.mdb
-rwxr-xr-x	root/root	usr/lib/beagle/beagled-index-helper
-rw-r--r--	root/root	usr/lib/beagle/libbeagleglue.a
-rwxr-xr-x	root/root	usr/lib/beagle/libbeagleglue.la
lrwxrwxrwx	root/root	usr/lib/beagle/libbeagleglue.so -> libbeagleglue.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/beagle/libbeagleglue.so.0 -> libbeagleglue.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/beagle/libbeagleglue.so.0.0.0
-rw-r--r--	root/root	usr/lib/libbeagle.a
-rwxr-xr-x	root/root	usr/lib/libbeagle.la
lrwxrwxrwx	root/root	usr/lib/libbeagle.so -> libbeagle.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libbeagle.so.0 -> libbeagle.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libbeagle.so.0.0.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/beagle-0.0.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/beagle-daemon.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/beagle-ui-0.0.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libbeagle-0.0.pc
drwxr-xr-x	root/root	usr/lib/python2.4/
drwxr-xr-x	root/root	usr/lib/python2.4/site-packages/
-rw-r--r--	root/root	usr/lib/python2.4/site-packages/beagle.a
-rwxr-xr-x	root/root	usr/lib/python2.4/site-packages/beagle.la
-rwxr-xr-x	root/root	usr/lib/python2.4/site-packages/beagle.so
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/beagle-config.1.gz
-rw-r--r--	root/root	usr/man/man1/beagle-query.1.gz
-rw-r--r--	root/root	usr/man/man1/beagle-shutdown.1.gz
-rw-r--r--	root/root	usr/man/man1/beagle-status.1.gz
-rw-r--r--	root/root	usr/man/man1/beagled.1.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/beagle-build-index.8.gz
-rw-r--r--	root/root	usr/man/man8/beagle-manage-index.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/beagle-build-index
-rwxr-xr-x	root/root	usr/sbin/beagle-manage-index
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/beagle/
