# yhafri/keepalived/.footprint
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/genhash
drwxr-xr-x	root/root	usr/etc/
drwxr-xr-x	root/root	usr/etc/init.d/
-rwxr-xr-x	root/root	usr/etc/init.d/keepalived
drwxr-xr-x	root/root	usr/etc/keepalived/
-rw-r--r--	root/root	usr/etc/keepalived/keepalived.conf
drwxr-xr-x	root/root	usr/etc/keepalived/samples/
-rw-r--r--	root/root	usr/etc/keepalived/samples/client.pem
-rw-r--r--	root/root	usr/etc/keepalived/samples/dh1024.pem
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.HTTP_GET.port
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.SMTP_CHECK
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.SSL_GET
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.fwmark
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.inhibit
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.misc_check
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.misc_check_arg
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.sample
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.status_code
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.track_interface
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.virtual_server_group
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.virtualhost
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.vrrp
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.vrrp.localcheck
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.vrrp.lvs_syncd
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.vrrp.routes
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.vrrp.scripts
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.vrrp.static_ipaddress
-rw-r--r--	root/root	usr/etc/keepalived/samples/keepalived.conf.vrrp.sync
-rw-r--r--	root/root	usr/etc/keepalived/samples/root.pem
-rw-r--r--	root/root	usr/etc/keepalived/samples/sample.misccheck.smbcheck.sh
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man/
drwxr-xr-x	root/root	usr/man/man/man1/
-rw-r--r--	root/root	usr/man/man/man1/genhash.1.gz
drwxr-xr-x	root/root	usr/man/man/man5/
-rw-r--r--	root/root	usr/man/man/man5/keepalived.conf.5.gz
drwxr-xr-x	root/root	usr/man/man/man8/
-rw-r--r--	root/root	usr/man/man/man8/keepalived.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwx------	root/root	usr/sbin/keepalived
