# opt/hal/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/hald
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/udev/
drwxr-xr-x	root/root	lib/udev/rules.d/
drwxr-xr-x	root/root	lib/udev/rules.d/rules.d/
-rw-r--r--	root/root	lib/udev/rules.d/rules.d/90-hal.rules
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/hal-device
-rwxr-xr-x	root/root	usr/bin/hal-disable-polling
-rwxr-xr-x	root/root	usr/bin/hal-find-by-capability
-rwxr-xr-x	root/root	usr/bin/hal-find-by-property
-rwxr-xr-x	root/root	usr/bin/hal-get-property
-rwxr-xr-x	root/root	usr/bin/hal-is-caller-locked-out
-rwxr-xr-x	root/root	usr/bin/hal-lock
-rwxr-xr-x	root/root	usr/bin/hal-set-property
-rwxr-xr-x	root/root	usr/bin/lshal
drwxr-xr-x	root/root	usr/etc/
drwxr-xr-x	root/root	usr/etc/dbus-1/
drwxr-xr-x	root/root	usr/etc/dbus-1/system.d/
-rw-r--r--	root/root	usr/etc/dbus-1/system.d/hal.conf
drwxr-xr-x	root/root	usr/etc/hal/
drwxr-xr-x	root/root	usr/etc/hal/fdi/
drwxr-xr-x	root/root	usr/etc/hal/fdi/information/
drwxr-xr-x	root/root	usr/etc/hal/fdi/policy/
drwxr-xr-x	root/root	usr/etc/hal/fdi/preprobe/
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/hal/
-rw-r--r--	root/root	usr/include/hal/libhal-storage.h
-rw-r--r--	root/root	usr/include/hal/libhal.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/hal/
-rwxr-xr-x	root/root	usr/lib/hal/hal-ipw-killswitch-linux
-rwxr-xr-x	root/root	usr/lib/hal/hal-storage-cleanup-all-mountpoints
-rwxr-xr-x	root/root	usr/lib/hal/hal-storage-cleanup-mountpoint
-rwxr-xr-x	root/root	usr/lib/hal/hal-storage-closetray
-rwxr-xr-x	root/root	usr/lib/hal/hal-storage-eject
-rwxr-xr-x	root/root	usr/lib/hal/hal-storage-mount
-rwxr-xr-x	root/root	usr/lib/hal/hal-storage-unmount
-rwxr-xr-x	root/root	usr/lib/hal/hal-system-power-pm-is-supported
-rwxr-xr-x	root/root	usr/lib/hal/hal-system-power-pmu
-rwxr-xr-x	root/root	usr/lib/hal/hal-system-setserial
-rwxr-xr-x	root/root	usr/lib/hal/hal-system-sonypic
-rwxr-xr-x	root/root	usr/lib/hal/hald-addon-acpi
-rwxr-xr-x	root/root	usr/lib/hal/hald-addon-cpufreq
-rwxr-xr-x	root/root	usr/lib/hal/hald-addon-hid-ups
-rwxr-xr-x	root/root	usr/lib/hal/hald-addon-imac-backlight
-rwxr-xr-x	root/root	usr/lib/hal/hald-addon-input
-rwxr-xr-x	root/root	usr/lib/hal/hald-addon-macbook-backlight
-rwxr-xr-x	root/root	usr/lib/hal/hald-addon-macbookpro-backlight
-rwxr-xr-x	root/root	usr/lib/hal/hald-addon-pmu
-rwxr-xr-x	root/root	usr/lib/hal/hald-addon-storage
-rwxr-xr-x	root/root	usr/lib/hal/hald-addon-usb-csr
-rwxr-xr-x	root/root	usr/lib/hal/hald-generate-fdi-cache
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-hiddev
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-ieee1394-unit
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-input
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-net-bluetooth
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-pc-floppy
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-printer
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-serial
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-smbios
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-storage
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-video4linux
-rwxr-xr-x	root/root	usr/lib/hal/hald-probe-volume
-rwxr-xr-x	root/root	usr/lib/hal/hald-runner
drwxr-xr-x	root/root	usr/lib/hal/scripts/
-rwxr-xr-x	root/root	usr/lib/hal/scripts/hal-dockstation-undock
-rwxr-xr-x	root/root	usr/lib/hal/scripts/hal-functions
-rwxr-xr-x	root/root	usr/lib/hal/scripts/hal-luks-remove
-rwxr-xr-x	root/root	usr/lib/hal/scripts/hal-luks-setup
-rwxr-xr-x	root/root	usr/lib/hal/scripts/hal-luks-teardown
-rwxr-xr-x	root/root	usr/lib/hal/scripts/hal-system-killswitch-get-power
-rwxr-xr-x	root/root	usr/lib/hal/scripts/hal-system-killswitch-set-power
-rwxr-xr-x	root/root	usr/lib/hal/scripts/hal-system-lcd-get-brightness
-rwxr-xr-x	root/root	usr/lib/hal/scripts/hal-system-lcd-set-brightness
-rwxr-xr-x	root/root	usr/lib/hal/scripts/hal-system-power-hibernate
-rwxr-xr-x	root/root	usr/lib/hal/scripts/hal-system-power-reboot
-rwxr-xr-x	root/root	usr/lib/hal/scripts/hal-system-power-set-power-save
-rwxr-xr-x	root/root	usr/lib/hal/scripts/hal-system-power-shutdown
-rwxr-xr-x	root/root	usr/lib/hal/scripts/hal-system-power-suspend
-rwxr-xr-x	root/root	usr/lib/hal/scripts/hal-system-power-suspend-hybrid
-rwxr-xr-x	root/root	usr/lib/hal/scripts/hal-system-wol-enable
-rwxr-xr-x	root/root	usr/lib/hal/scripts/hal-system-wol-enabled
-rwxr-xr-x	root/root	usr/lib/hal/scripts/hal-system-wol-supported
drwxr-xr-x	root/root	usr/lib/hal/scripts/linux/
-rwxr-xr-x	root/root	usr/lib/hal/scripts/linux/hal-dockstation-undock-linux
-rwxr-xr-x	root/root	usr/lib/hal/scripts/linux/hal-luks-remove-linux
-rwxr-xr-x	root/root	usr/lib/hal/scripts/linux/hal-luks-setup-linux
-rwxr-xr-x	root/root	usr/lib/hal/scripts/linux/hal-luks-teardown-linux
-rwxr-xr-x	root/root	usr/lib/hal/scripts/linux/hal-system-killswitch-get-power-linux
-rwxr-xr-x	root/root	usr/lib/hal/scripts/linux/hal-system-killswitch-set-power-linux
-rwxr-xr-x	root/root	usr/lib/hal/scripts/linux/hal-system-lcd-get-brightness-linux
-rwxr-xr-x	root/root	usr/lib/hal/scripts/linux/hal-system-lcd-set-brightness-linux
-rwxr-xr-x	root/root	usr/lib/hal/scripts/linux/hal-system-power-hibernate-linux
-rwxr-xr-x	root/root	usr/lib/hal/scripts/linux/hal-system-power-reboot-linux
-rwxr-xr-x	root/root	usr/lib/hal/scripts/linux/hal-system-power-set-power-save-linux
-rwxr-xr-x	root/root	usr/lib/hal/scripts/linux/hal-system-power-shutdown-linux
-rwxr-xr-x	root/root	usr/lib/hal/scripts/linux/hal-system-power-suspend-hybrid-linux
-rwxr-xr-x	root/root	usr/lib/hal/scripts/linux/hal-system-power-suspend-linux
lrwxrwxrwx	root/root	usr/lib/hal/scripts/linux/hal-system-wol-enable-linux -> hal-system-wol-linux
lrwxrwxrwx	root/root	usr/lib/hal/scripts/linux/hal-system-wol-enabled-linux -> hal-system-wol-linux
-rwxr-xr-x	root/root	usr/lib/hal/scripts/linux/hal-system-wol-linux
lrwxrwxrwx	root/root	usr/lib/hal/scripts/linux/hal-system-wol-supported-linux -> hal-system-wol-linux
-rw-r--r--	root/root	usr/lib/libhal-storage.a
-rwxr-xr-x	root/root	usr/lib/libhal-storage.la
lrwxrwxrwx	root/root	usr/lib/libhal-storage.so -> libhal-storage.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libhal-storage.so.1 -> libhal-storage.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/libhal-storage.so.1.0.0
-rw-r--r--	root/root	usr/lib/libhal.a
-rwxr-xr-x	root/root	usr/lib/libhal.la
lrwxrwxrwx	root/root	usr/lib/libhal.so -> libhal.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libhal.so.1 -> libhal.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/libhal.so.1.0.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/hal-storage.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/hal.pc
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/hal-disable-polling.1.gz
-rw-r--r--	root/root	usr/man/man1/hal-find-by-capability.1.gz
-rw-r--r--	root/root	usr/man/man1/hal-find-by-property.1.gz
-rw-r--r--	root/root	usr/man/man1/hal-get-property.1.gz
-rw-r--r--	root/root	usr/man/man1/hal-is-caller-locked-out.1.gz
-rw-r--r--	root/root	usr/man/man1/hal-lock.1.gz
-rw-r--r--	root/root	usr/man/man1/hal-set-property.1.gz
-rw-r--r--	root/root	usr/man/man1/lshal.1.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/hald.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/hald
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/hal/
drwxr-xr-x	root/root	usr/share/hal/fdi/
-rw-r--r--	root/root	usr/share/hal/fdi/fdi.dtd
drwxr-xr-x	root/root	usr/share/hal/fdi/information/
drwxr-xr-x	root/root	usr/share/hal/fdi/information/10freedesktop/
drwxr-xr-x	root/root	usr/share/hal/fdi/information/20thirdparty/
drwxr-xr-x	root/root	usr/share/hal/fdi/policy/
drwxr-xr-x	root/root	usr/share/hal/fdi/policy/10osvendor/
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/10-cpufreq.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/10-dockstation.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/10-imac-backlight.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/10-input-policy.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/10-laptop-panel-mgmt-policy.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/10-macbook-backlight.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/10-macbookpro-utils.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/10-power-mgmt-policy.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/10-rfkill-switch.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/10-tabletPCs.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/10-usbcsr-mice.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/10-x11-input.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/15-storage-luks.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/20-storage-methods.fdi
-rw-r--r--	root/root	usr/share/hal/fdi/policy/10osvendor/30-wol.fdi
drwxr-xr-x	root/root	usr/share/hal/fdi/policy/20thirdparty/
drwxr-xr-x	root/root	usr/share/hal/fdi/preprobe/
drwxr-xr-x	root/root	usr/share/hal/fdi/preprobe/10osvendor/
drwxr-xr-x	root/root	usr/share/hal/fdi/preprobe/20thirdparty/
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/cache/
drwxr-xr-x	haldaemon/haldaemon	var/cache/hald/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	haldaemon/haldaemon	var/lib/hal/
drwxr-xr-x	root/root	var/run/
drwxr-xr-x	haldaemon/haldaemon	var/run/hald/
drwxr-xr-x	root/root	var/run/hald/hald-local/
drwxr-xr-x	root/root	var/run/hald/hald-runner/
