# yhafri/alsa-driver/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/sound
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/modules/
drwxr-xr-x	root/root	lib/modules/<kernel-version>/
drwxr-xr-x	root/root	lib/modules/<kernel-version>/kernel/
drwxr-xr-x	root/root	lib/modules/<kernel-version>/kernel/sound/
drwxr-xr-x	root/root	lib/modules/<kernel-version>/kernel/sound/acore/
drwxr-xr-x	root/root	lib/modules/<kernel-version>/kernel/sound/acore/oss/
-rw-r--r--	root/root	lib/modules/<kernel-version>/kernel/sound/acore/oss/snd-mixer-oss.ko
-rw-r--r--	root/root	lib/modules/<kernel-version>/kernel/sound/acore/oss/snd-pcm-oss.ko
drwxr-xr-x	root/root	lib/modules/<kernel-version>/kernel/sound/acore/seq/
drwxr-xr-x	root/root	lib/modules/<kernel-version>/kernel/sound/acore/seq/oss/
-rw-r--r--	root/root	lib/modules/<kernel-version>/kernel/sound/acore/seq/oss/snd-seq-oss.ko
-rw-r--r--	root/root	lib/modules/<kernel-version>/kernel/sound/acore/seq/snd-seq-device.ko
-rw-r--r--	root/root	lib/modules/<kernel-version>/kernel/sound/acore/seq/snd-seq-midi-event.ko
-rw-r--r--	root/root	lib/modules/<kernel-version>/kernel/sound/acore/seq/snd-seq.ko
-rw-r--r--	root/root	lib/modules/<kernel-version>/kernel/sound/acore/snd-page-alloc.ko
-rw-r--r--	root/root	lib/modules/<kernel-version>/kernel/sound/acore/snd-pcm.ko
-rw-r--r--	root/root	lib/modules/<kernel-version>/kernel/sound/acore/snd-timer.ko
-rw-r--r--	root/root	lib/modules/<kernel-version>/kernel/sound/acore/snd.ko
drwxr-xr-x	root/root	lib/modules/<kernel-version>/kernel/sound/pci/
drwxr-xr-x	root/root	lib/modules/<kernel-version>/kernel/sound/pci/ac97/
-rw-r--r--	root/root	lib/modules/<kernel-version>/kernel/sound/pci/ac97/snd-ac97-bus.ko
-rw-r--r--	root/root	lib/modules/<kernel-version>/kernel/sound/pci/ac97/snd-ac97-codec.ko
-rw-r--r--	root/root	lib/modules/<kernel-version>/kernel/sound/pci/snd-intel8x0.ko
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/sound/
-rw-r--r--	root/root	usr/include/sound/ac97_codec.h
-rw-r--r--	root/root	usr/include/sound/ad1816a.h
-rw-r--r--	root/root	usr/include/sound/ad1848.h
-rw-r--r--	root/root	usr/include/sound/ainstr_fm.h
-rw-r--r--	root/root	usr/include/sound/ainstr_gf1.h
-rw-r--r--	root/root	usr/include/sound/ainstr_iw.h
-rw-r--r--	root/root	usr/include/sound/ainstr_simple.h
-rw-r--r--	root/root	usr/include/sound/ak4114.h
-rw-r--r--	root/root	usr/include/sound/ak4117.h
-rw-r--r--	root/root	usr/include/sound/ak4531_codec.h
-rw-r--r--	root/root	usr/include/sound/ak4xxx-adda.h
-rw-r--r--	root/root	usr/include/sound/asequencer.h
-rw-r--r--	root/root	usr/include/sound/asound.h
-rw-r--r--	root/root	usr/include/sound/asound_fm.h
-rw-r--r--	root/root	usr/include/sound/asoundef.h
-rw-r--r--	root/root	usr/include/sound/control.h
-rw-r--r--	root/root	usr/include/sound/core.h
-rw-r--r--	root/root	usr/include/sound/cs4231.h
-rw-r--r--	root/root	usr/include/sound/cs46xx.h
-rw-r--r--	root/root	usr/include/sound/cs46xx_dsp_scb_types.h
-rw-r--r--	root/root	usr/include/sound/cs46xx_dsp_spos.h
-rw-r--r--	root/root	usr/include/sound/cs46xx_dsp_task_types.h
-rw-r--r--	root/root	usr/include/sound/cs8403.h
-rw-r--r--	root/root	usr/include/sound/cs8427.h
-rw-r--r--	root/root	usr/include/sound/driver.h
-rw-r--r--	root/root	usr/include/sound/emu10k1.h
-rw-r--r--	root/root	usr/include/sound/emu10k1_synth.h
-rw-r--r--	root/root	usr/include/sound/emu8000.h
-rw-r--r--	root/root	usr/include/sound/emu8000_reg.h
-rw-r--r--	root/root	usr/include/sound/emux_legacy.h
-rw-r--r--	root/root	usr/include/sound/emux_synth.h
-rw-r--r--	root/root	usr/include/sound/es1688.h
-rw-r--r--	root/root	usr/include/sound/gus.h
-rw-r--r--	root/root	usr/include/sound/hdsp.h
-rw-r--r--	root/root	usr/include/sound/hdspm.h
-rw-r--r--	root/root	usr/include/sound/hwdep.h
-rw-r--r--	root/root	usr/include/sound/i2c.h
-rw-r--r--	root/root	usr/include/sound/info.h
-rw-r--r--	root/root	usr/include/sound/initval.h
-rw-r--r--	root/root	usr/include/sound/memalloc.h
-rw-r--r--	root/root	usr/include/sound/minors.h
-rw-r--r--	root/root	usr/include/sound/mixer_oss.h
-rw-r--r--	root/root	usr/include/sound/mpu401.h
-rw-r--r--	root/root	usr/include/sound/opl3.h
-rw-r--r--	root/root	usr/include/sound/opl4.h
-rw-r--r--	root/root	usr/include/sound/pcm-indirect.h
-rw-r--r--	root/root	usr/include/sound/pcm.h
-rw-r--r--	root/root	usr/include/sound/pcm_oss.h
-rw-r--r--	root/root	usr/include/sound/pcm_params.h
-rw-r--r--	root/root	usr/include/sound/rawmidi.h
-rw-r--r--	root/root	usr/include/sound/sb.h
-rw-r--r--	root/root	usr/include/sound/sb16_csp.h
-rw-r--r--	root/root	usr/include/sound/seq_device.h
-rw-r--r--	root/root	usr/include/sound/seq_instr.h
-rw-r--r--	root/root	usr/include/sound/seq_kernel.h
-rw-r--r--	root/root	usr/include/sound/seq_midi_emul.h
-rw-r--r--	root/root	usr/include/sound/seq_midi_event.h
-rw-r--r--	root/root	usr/include/sound/seq_oss.h
-rw-r--r--	root/root	usr/include/sound/seq_oss_legacy.h
-rw-r--r--	root/root	usr/include/sound/seq_virmidi.h
-rw-r--r--	root/root	usr/include/sound/sfnt_info.h
-rw-r--r--	root/root	usr/include/sound/snd_wavefront.h
-rw-r--r--	root/root	usr/include/sound/soundfont.h
-rw-r--r--	root/root	usr/include/sound/sscape_ioctl.h
-rw-r--r--	root/root	usr/include/sound/tea575x-tuner.h
-rw-r--r--	root/root	usr/include/sound/tea6330t.h
-rw-r--r--	root/root	usr/include/sound/timer.h
-rw-r--r--	root/root	usr/include/sound/tlv.h
-rw-r--r--	root/root	usr/include/sound/trident.h
-rw-r--r--	root/root	usr/include/sound/typedefs.h
-rw-r--r--	root/root	usr/include/sound/uda1341.h
-rw-r--r--	root/root	usr/include/sound/util_mem.h
-rw-r--r--	root/root	usr/include/sound/version.h
-rw-r--r--	root/root	usr/include/sound/vx_core.h
-rw-r--r--	root/root	usr/include/sound/wavefront.h
-rw-r--r--	root/root	usr/include/sound/wavefront_fx.h
-rw-r--r--	root/root	usr/include/sound/ymfpci.h
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/snddevices
