# yhafri/xinc/.footprint
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/xinc
drwxr-xr-x	root/root	usr/etc/
drwxr-xr-x	root/root	usr/etc/init.d/
-rwxr-xr-x	root/root	usr/etc/init.d/xinc
drwxr-xr-x	root/root	usr/etc/xinc/
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/php/
-rw-r--r--	root/root	usr/share/php/Xinc.php
drwxr-xr-x	root/root	usr/share/php/Xinc/
drwxr-xr-x	root/root	usr/share/php/Xinc/Builder/
-rw-r--r--	root/root	usr/share/php/Xinc/Builder/Interface.php
-rw-r--r--	root/root	usr/share/php/Xinc/Builder/Phing.php
drwxr-xr-x	root/root	usr/share/php/Xinc/Exception/
-rw-r--r--	root/root	usr/share/php/Xinc/Exception/MalformedConfig.php
-rw-r--r--	root/root	usr/share/php/Xinc/Exception/ModificationSet.php
-rw-r--r--	root/root	usr/share/php/Xinc/Logger.php
drwxr-xr-x	root/root	usr/share/php/Xinc/Logger/
-rw-r--r--	root/root	usr/share/php/Xinc/Logger/Message.php
drwxr-xr-x	root/root	usr/share/php/Xinc/ModificationSet/
-rw-r--r--	root/root	usr/share/php/Xinc/ModificationSet/Fake.php
-rw-r--r--	root/root	usr/share/php/Xinc/ModificationSet/Interface.php
-rw-r--r--	root/root	usr/share/php/Xinc/ModificationSet/Svn.php
-rw-r--r--	root/root	usr/share/php/Xinc/Parser.php
-rw-r--r--	root/root	usr/share/php/Xinc/Project.php
drwxr-xr-x	root/root	usr/share/php/Xinc/Publisher/
-rw-r--r--	root/root	usr/share/php/Xinc/Publisher/Email.php
-rw-r--r--	root/root	usr/share/php/Xinc/Publisher/Fake.php
-rw-r--r--	root/root	usr/share/php/Xinc/Publisher/Interface.php
-rw-r--r--	root/root	usr/share/php/Xinc/Publisher/Phing.php
