# yhafri/monit/.footprint
drwxr-xr-x	root/root	etc/
-rwx------	root/root	etc/monitrc
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/monit
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-r-xr-xr-x	root/root	usr/bin/monit
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-r--r--r--	root/root	usr/man/man1/monit.1.gz
