# attic/netkit-rwho/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/rwhod
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/ruptime
-rwxr-xr-x	root/root	usr/bin/rwho
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/ruptime.1.gz
-rw-r--r--	root/root	usr/man/man1/rwho.1.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/rwhod.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/rwhod
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/spool/
drwxr-xr-x	root/root	var/spool/rwho/
