# yhafri/libmysqld/.footprint
drwxr-xr-x	root/root	etc/
-rw-------	root/root	etc/my.cnf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/mysqld
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/comp_err
-rwxr-xr-x	root/root	usr/bin/innochecksum
-rwxr-xr-x	root/root	usr/bin/msql2mysql
-rwxr-xr-x	root/root	usr/bin/my_print_defaults
-rwxr-xr-x	root/root	usr/bin/myisam_ftdump
-rwxr-xr-x	root/root	usr/bin/myisamchk
-rwxr-xr-x	root/root	usr/bin/myisamlog
-rwxr-xr-x	root/root	usr/bin/myisampack
-rwxr-xr-x	root/root	usr/bin/mysql
-rwxr-xr-x	root/root	usr/bin/mysql_client_test
-rwxr-xr-x	root/root	usr/bin/mysql_config
-rwxr-xr-x	root/root	usr/bin/mysql_convert_table_format
-rwxr-xr-x	root/root	usr/bin/mysql_create_system_tables
-rwxr-xr-x	root/root	usr/bin/mysql_explain_log
-rwxr-xr-x	root/root	usr/bin/mysql_find_rows
-rwxr-xr-x	root/root	usr/bin/mysql_fix_extensions
-rwxr-xr-x	root/root	usr/bin/mysql_fix_privilege_tables
-rwxr-xr-x	root/root	usr/bin/mysql_install_db
-rwxr-xr-x	root/root	usr/bin/mysql_secure_installation
-rwxr-xr-x	root/root	usr/bin/mysql_setpermission
-rwxr-xr-x	root/root	usr/bin/mysql_tzinfo_to_sql
-rwxr-xr-x	root/root	usr/bin/mysql_upgrade
-rwxr-xr-x	root/root	usr/bin/mysql_waitpid
-rwxr-xr-x	root/root	usr/bin/mysql_zap
-rwxr-xr-x	root/root	usr/bin/mysqlaccess
-rwxr-xr-x	root/root	usr/bin/mysqladmin
-rwxr-xr-x	root/root	usr/bin/mysqlbinlog
-rwxr-xr-x	root/root	usr/bin/mysqlbug
-rwxr-xr-x	root/root	usr/bin/mysqlcheck
-rwxr-xr-x	root/root	usr/bin/mysqld_multi
-rwxr-xr-x	root/root	usr/bin/mysqld_safe
-rwxr-xr-x	root/root	usr/bin/mysqldump
-rwxr-xr-x	root/root	usr/bin/mysqldumpslow
-rwxr-xr-x	root/root	usr/bin/mysqlhotcopy
-rwxr-xr-x	root/root	usr/bin/mysqlimport
-rwxr-xr-x	root/root	usr/bin/mysqlshow
-rwxr-xr-x	root/root	usr/bin/mysqlslap
-rwxr-xr-x	root/root	usr/bin/mysqltest
-rwxr-xr-x	root/root	usr/bin/mysqltestmanager
-rwxr-xr-x	root/root	usr/bin/mysqltestmanager-pwgen
-rwxr-xr-x	root/root	usr/bin/mysqltestmanagerc
-rwxr-xr-x	root/root	usr/bin/perror
-rwxr-xr-x	root/root	usr/bin/replace
-rwxr-xr-x	root/root	usr/bin/resolve_stack_dump
-rwxr-xr-x	root/root	usr/bin/resolveip
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/mysql/
-rw-r--r--	root/root	usr/include/mysql/errmsg.h
-rw-r--r--	root/root	usr/include/mysql/keycache.h
-rw-r--r--	root/root	usr/include/mysql/m_ctype.h
-rw-r--r--	root/root	usr/include/mysql/m_string.h
-rw-r--r--	root/root	usr/include/mysql/my_alloc.h
-rw-r--r--	root/root	usr/include/mysql/my_config.h
-rw-r--r--	root/root	usr/include/mysql/my_dbug.h
-rw-r--r--	root/root	usr/include/mysql/my_dir.h
-rw-r--r--	root/root	usr/include/mysql/my_getopt.h
-rw-r--r--	root/root	usr/include/mysql/my_global.h
-rw-r--r--	root/root	usr/include/mysql/my_list.h
-rw-r--r--	root/root	usr/include/mysql/my_net.h
-rw-r--r--	root/root	usr/include/mysql/my_no_pthread.h
-rw-r--r--	root/root	usr/include/mysql/my_pthread.h
-rw-r--r--	root/root	usr/include/mysql/my_semaphore.h
-rw-r--r--	root/root	usr/include/mysql/my_sys.h
-rw-r--r--	root/root	usr/include/mysql/my_xml.h
-rw-r--r--	root/root	usr/include/mysql/mysql.h
-rw-r--r--	root/root	usr/include/mysql/mysql_com.h
-rw-r--r--	root/root	usr/include/mysql/mysql_embed.h
-rw-r--r--	root/root	usr/include/mysql/mysql_time.h
-rw-r--r--	root/root	usr/include/mysql/mysql_version.h
-rw-r--r--	root/root	usr/include/mysql/mysqld_ername.h
-rw-r--r--	root/root	usr/include/mysql/mysqld_error.h
-rw-r--r--	root/root	usr/include/mysql/plugin.h
-rw-r--r--	root/root	usr/include/mysql/raid.h
-rw-r--r--	root/root	usr/include/mysql/readline.h
-rw-r--r--	root/root	usr/include/mysql/sql_common.h
-rw-r--r--	root/root	usr/include/mysql/sql_state.h
-rw-r--r--	root/root	usr/include/mysql/sslopt-case.h
-rw-r--r--	root/root	usr/include/mysql/sslopt-longopts.h
-rw-r--r--	root/root	usr/include/mysql/sslopt-vars.h
-rw-r--r--	root/root	usr/include/mysql/typelib.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/mysql/
-rw-r--r--	root/root	usr/lib/mysql/libcsv.a
-rw-r--r--	root/root	usr/lib/mysql/libdbug.a
-rw-r--r--	root/root	usr/lib/mysql/libheap.a
-rw-r--r--	root/root	usr/lib/mysql/libmyisam.a
-rw-r--r--	root/root	usr/lib/mysql/libmyisammrg.a
-rw-r--r--	root/root	usr/lib/mysql/libmysqlclient.a
-rwxr-xr-x	root/root	usr/lib/mysql/libmysqlclient.la
lrwxrwxrwx	root/root	usr/lib/mysql/libmysqlclient.so -> libmysqlclient.so.15.0.0
lrwxrwxrwx	root/root	usr/lib/mysql/libmysqlclient.so.15 -> libmysqlclient.so.15.0.0
-rwxr-xr-x	root/root	usr/lib/mysql/libmysqlclient.so.15.0.0
-rw-r--r--	root/root	usr/lib/mysql/libmysqlclient_r.a
-rwxr-xr-x	root/root	usr/lib/mysql/libmysqlclient_r.la
lrwxrwxrwx	root/root	usr/lib/mysql/libmysqlclient_r.so -> libmysqlclient_r.so.15.0.0
lrwxrwxrwx	root/root	usr/lib/mysql/libmysqlclient_r.so.15 -> libmysqlclient_r.so.15.0.0
-rwxr-xr-x	root/root	usr/lib/mysql/libmysqlclient_r.so.15.0.0
-rw-r--r--	root/root	usr/lib/mysql/libmystrings.a
-rw-r--r--	root/root	usr/lib/mysql/libmysys.a
-rw-r--r--	root/root	usr/lib/mysql/libvio.a
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/msql2mysql.1.gz
-rw-r--r--	root/root	usr/man/man1/myisamchk.1.gz
-rw-r--r--	root/root	usr/man/man1/myisamlog.1.gz
-rw-r--r--	root/root	usr/man/man1/myisampack.1.gz
-rw-r--r--	root/root	usr/man/man1/mysql.1.gz
-rw-r--r--	root/root	usr/man/man1/mysql.server.1.gz
-rw-r--r--	root/root	usr/man/man1/mysql_config.1.gz
-rw-r--r--	root/root	usr/man/man1/mysql_fix_privilege_tables.1.gz
-rw-r--r--	root/root	usr/man/man1/mysql_upgrade.1.gz
-rw-r--r--	root/root	usr/man/man1/mysql_zap.1.gz
-rw-r--r--	root/root	usr/man/man1/mysqlaccess.1.gz
-rw-r--r--	root/root	usr/man/man1/mysqladmin.1.gz
-rw-r--r--	root/root	usr/man/man1/mysqlbinlog.1.gz
-rw-r--r--	root/root	usr/man/man1/mysqlcheck.1.gz
-rw-r--r--	root/root	usr/man/man1/mysqld.1.gz
-rw-r--r--	root/root	usr/man/man1/mysqld_multi.1.gz
-rw-r--r--	root/root	usr/man/man1/mysqld_safe.1.gz
-rw-r--r--	root/root	usr/man/man1/mysqldump.1.gz
-rw-r--r--	root/root	usr/man/man1/mysqlhotcopy.1.gz
-rw-r--r--	root/root	usr/man/man1/mysqlimport.1.gz
-rw-r--r--	root/root	usr/man/man1/mysqlman.1.gz
-rw-r--r--	root/root	usr/man/man1/mysqlmanager.1.gz
-rw-r--r--	root/root	usr/man/man1/mysqlshow.1.gz
-rw-r--r--	root/root	usr/man/man1/mysqlslap.1.gz
-rw-r--r--	root/root	usr/man/man1/perror.1.gz
-rw-r--r--	root/root	usr/man/man1/replace.1.gz
-rw-r--r--	root/root	usr/man/man1/safe_mysqld.1.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/mysqld
-rwxr-xr-x	root/root	usr/sbin/mysqlmanager
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/mysql/
drwxr-xr-x	root/root	usr/share/mysql/charsets/
-rw-r--r--	root/root	usr/share/mysql/charsets/Index.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/README
-rw-r--r--	root/root	usr/share/mysql/charsets/armscii8.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/ascii.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/cp1250.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/cp1251.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/cp1256.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/cp1257.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/cp850.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/cp852.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/cp866.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/dec8.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/geostd8.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/greek.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/hebrew.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/hp8.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/keybcs2.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/koi8r.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/koi8u.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/latin1.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/latin2.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/latin5.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/latin7.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/macce.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/macroman.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/swe7.xml
drwxr-xr-x	root/root	usr/share/mysql/english/
-rw-r--r--	root/root	usr/share/mysql/english/errmsg.sys
-rw-r--r--	root/root	usr/share/mysql/fill_help_tables.sql
-rw-r--r--	root/root	usr/share/mysql/my-huge.cnf
-rw-r--r--	root/root	usr/share/mysql/my-innodb-heavy-4G.cnf
-rw-r--r--	root/root	usr/share/mysql/my-large.cnf
-rw-r--r--	root/root	usr/share/mysql/my-medium.cnf
-rw-r--r--	root/root	usr/share/mysql/my-small.cnf
-rw-r--r--	root/root	usr/share/mysql/mysql_fix_privilege_tables.sql
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwx------	root/root	var/lib/mysql/
drwxr-xr-x	root/root	var/log/
-rw-------	root/root	var/log/mysqld.log (EMPTY)
