# contrib-romster/distcc/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/distccd
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/distcc
-rwxr-xr-x	root/root	usr/bin/distccd
-rwxr-xr-x	root/root	usr/bin/distccmon-text
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/distcc.1.gz
-rw-r--r--	root/root	usr/man/man1/distccd.1.gz
-rw-r--r--	root/root	usr/man/man1/distccmon-text.1.gz
