# ports/privoxy/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/privoxy/
-rw-r--r--	root/root	etc/privoxy/config
-rw-r--r--	root/root	etc/privoxy/default.action
-rw-r--r--	root/root	etc/privoxy/default.filter
-rw-r--r--	root/root	etc/privoxy/match-all.action
-rw-r--r--	root/root	etc/privoxy/regression-tests.action
drwxr-xr-x	root/root	etc/privoxy/templates/
-rw-r--r--	root/root	etc/privoxy/templates/blocked
-rw-r--r--	root/root	etc/privoxy/templates/cgi-error-404
-rw-r--r--	root/root	etc/privoxy/templates/cgi-error-bad-param
-rw-r--r--	root/root	etc/privoxy/templates/cgi-error-disabled
-rw-r--r--	root/root	etc/privoxy/templates/cgi-error-file
-rw-r--r--	root/root	etc/privoxy/templates/cgi-error-file-read-only
-rw-r--r--	root/root	etc/privoxy/templates/cgi-error-modified
-rw-r--r--	root/root	etc/privoxy/templates/cgi-error-parse
-rw-r--r--	root/root	etc/privoxy/templates/cgi-style.css
-rw-r--r--	root/root	etc/privoxy/templates/connect-failed
-rw-r--r--	root/root	etc/privoxy/templates/connection-timeout
-rw-r--r--	root/root	etc/privoxy/templates/default
-rw-r--r--	root/root	etc/privoxy/templates/edit-actions-add-url-form
-rw-r--r--	root/root	etc/privoxy/templates/edit-actions-for-url
-rw-r--r--	root/root	etc/privoxy/templates/edit-actions-for-url-filter
-rw-r--r--	root/root	etc/privoxy/templates/edit-actions-list
-rw-r--r--	root/root	etc/privoxy/templates/edit-actions-list-button
-rw-r--r--	root/root	etc/privoxy/templates/edit-actions-list-section
-rw-r--r--	root/root	etc/privoxy/templates/edit-actions-list-url
-rw-r--r--	root/root	etc/privoxy/templates/edit-actions-remove-url-form
-rw-r--r--	root/root	etc/privoxy/templates/edit-actions-url-form
-rw-r--r--	root/root	etc/privoxy/templates/forwarding-failed
-rw-r--r--	root/root	etc/privoxy/templates/mod-local-help
-rw-r--r--	root/root	etc/privoxy/templates/mod-support-and-service
-rw-r--r--	root/root	etc/privoxy/templates/mod-title
-rw-r--r--	root/root	etc/privoxy/templates/mod-unstable-warning
-rw-r--r--	root/root	etc/privoxy/templates/no-server-data
-rw-r--r--	root/root	etc/privoxy/templates/no-such-domain
-rw-r--r--	root/root	etc/privoxy/templates/show-request
-rw-r--r--	root/root	etc/privoxy/templates/show-status
-rw-r--r--	root/root	etc/privoxy/templates/show-status-file
-rw-r--r--	root/root	etc/privoxy/templates/show-url-info
-rw-r--r--	root/root	etc/privoxy/templates/show-version
-rw-r--r--	root/root	etc/privoxy/templates/toggle
-rw-r--r--	root/root	etc/privoxy/templates/toggle-mini
-rw-r--r--	root/root	etc/privoxy/templates/untrusted
-rw-r--r--	root/root	etc/privoxy/templates/url-info-osd.xml
-rw-r--r--	root/root	etc/privoxy/user.action
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/privoxy
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/privoxy.1.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/privoxy
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
drwxr-xr-x	root/root	var/log/privoxy/
