# yhafri/smartmontools/.footprint
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/etc/
drwxr-xr-x	root/root	usr/etc/rc.d/
drwxr-xr-x	root/root	usr/etc/rc.d/init.d/
-rwxr-xr-x	root/root	usr/etc/rc.d/init.d/smartd
-rw-r--r--	root/root	usr/etc/smartd.conf
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man5/
-rw-r--r--	root/root	usr/man/man5/smartd.conf.5.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/smartctl.8.gz
-rw-r--r--	root/root	usr/man/man8/smartd.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/smartctl
-rwxr-xr-x	root/root	usr/sbin/smartd
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/doc/
drwxr-xr-x	root/root	usr/share/doc/smartmontools-5.33/
-rw-r--r--	root/root	usr/share/doc/smartmontools-5.33/AUTHORS
-rw-r--r--	root/root	usr/share/doc/smartmontools-5.33/CHANGELOG
-rw-r--r--	root/root	usr/share/doc/smartmontools-5.33/COPYING
-rw-r--r--	root/root	usr/share/doc/smartmontools-5.33/INSTALL
-rw-r--r--	root/root	usr/share/doc/smartmontools-5.33/NEWS
-rw-r--r--	root/root	usr/share/doc/smartmontools-5.33/README
-rw-r--r--	root/root	usr/share/doc/smartmontools-5.33/TODO
-rw-r--r--	root/root	usr/share/doc/smartmontools-5.33/WARNINGS
drwxr-xr-x	root/root	usr/share/doc/smartmontools-5.33/examplescripts/
-rwxr-xr-x	root/root	usr/share/doc/smartmontools-5.33/examplescripts/Example1
-rwxr-xr-x	root/root	usr/share/doc/smartmontools-5.33/examplescripts/Example2
-rwxr-xr-x	root/root	usr/share/doc/smartmontools-5.33/examplescripts/Example3
-rw-r--r--	root/root	usr/share/doc/smartmontools-5.33/examplescripts/README
-rw-r--r--	root/root	usr/share/doc/smartmontools-5.33/smartd.conf
