# denyhosts/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/denyhosts/
-rw-r--r--	root/root	etc/denyhosts/denyhosts.conf
drwxr-xr-x	root/root	etc/rc.d/
lrwxrwxrwx	root/root	etc/rc.d/denyhosts -> /usr/share/denyhosts/daemon-control
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/denyhosts.py
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/python2.5/
drwxr-xr-x	root/root	usr/lib/python2.5/site-packages/
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts-2.5-py2.5.egg-info
drwxr-xr-x	root/root	usr/lib/python2.5/site-packages/DenyHosts/
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/__init__.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/__init__.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/allowedhosts.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/allowedhosts.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/constants.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/constants.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/counter.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/counter.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/daemon.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/daemon.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/deny_hosts.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/deny_hosts.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/denyfileutil.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/denyfileutil.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/filetracker.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/filetracker.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/lockfile.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/lockfile.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/loginattempt.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/loginattempt.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/old-daemon.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/old-daemon.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/plugin.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/plugin.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/prefs.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/prefs.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/purgecounter.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/purgecounter.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/python_version.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/python_version.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/regex.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/regex.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/report.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/report.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/restricted.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/restricted.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/sync.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/sync.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/util.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/util.pyc
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/version.py
-rw-r--r--	root/root	usr/lib/python2.5/site-packages/DenyHosts/version.pyc
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/denyhosts/
-rwx------	root/root	usr/share/denyhosts/daemon-control
drwxr-xr-x	root/root	usr/share/denyhosts/plugins/
-rw-r--r--	root/root	usr/share/denyhosts/plugins/README.contrib
-rw-r--r--	root/root	usr/share/denyhosts/plugins/shorewall_allow.sh
-rw-r--r--	root/root	usr/share/denyhosts/plugins/shorewall_deny.sh
-rwxr-xr-x	root/root	usr/share/denyhosts/plugins/test_deny.py
drwxr-xr-x	root/root	usr/share/denyhosts/scripts/
-rwxr-xr-x	root/root	usr/share/denyhosts/scripts/restricted_from_invalid.py
-rwxr-xr-x	root/root	usr/share/denyhosts/scripts/restricted_from_passwd.py
-rw-r--r--	root/root	usr/share/denyhosts/setup.py
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
